# supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd

[program:blitz]
command=/bin/bash /home/blitz/blitz/blitz.sh
user=blitz
directory=/home/blitz/blitz
